# ansible/roles/prod/templates/teleplata.service.j2
[Unit]
Description=TeleplataService

[Service]
Type=idle
ExecStart={{ ENV_PATH }}/bin/gunicorn teleplata.main:"create_app()" -w 4
Restart=always

[Install]
WantedBy=multi-user.target